# systemd-units/blfs/units/nfs-client.target
[Unit]
Description=NFS Client Daemons
Before=remote-fs-pre.target
Wants=rpc-statd.service

[Install]
WantedBy=multi-user.target
